# ff.conf
[program:ff]
command=python3 /home/curator2/ff/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ff/test.err.log
stdout_logfile=/home/curator2/ff/test.out.log